# daemon/kdeconnectd.service.in
# SPDX-FileCopyrightText: 2026 Nicolas Fella <nicolas.fella@gmx.de>
# SPDX-License-Identifier: CC0-1.0

[Unit]
Description=KDE Connect Daemon
PartOf=graphical-session.target
After=graphical-session.target

[Service]
ExecStart=@KDE_INSTALL_FULL_BINDIR@/kdeconnectd
Type=exec
BusName=org.kde.kdeconnectd
Slice=background.slice
Restart=on-failure

[Install]
WantedBy=graphical-session.target
